# upd72020x-load-openrc
#!/sbin/openrc-run

start() {
	/usr/bin/upd72020x-check-and-init
}

